# images/3.3.1-wp4.8.2-php7.2/scripts/supervisord-zwoocommerce.conf
[program:woocommerce-setup]
command=sh -x /setup-woocommerce.sh
startsecs=1
autorestart=false
startretries=1